# images/base/files/etc/systemd/system/kubelet.slice
[Unit]
Description=slice used to run Kubernetes / Kubelet
Before=slices.target

[Slice]
MemoryAccounting=true
CPUAccounting=true
